# zkillquery.container
[Unit]
Description=ZKillQuery Container

[Container]
Image=localhost/zkillquery:latest
ContainerName=zkillquery-run
Volume=~/ZKillQueryData:/app/ZKillQueryData
Exec=/bin/bash -c "cd /app/ZKillQuery && ls -ld * && python3 -uZKillQuery.py"

[Service]
Restart=always  # Auto-restart if it exits/crashes

[Install]
WantedBy=default.target

